# ansible/ohsu-migration/roles/monarch-noctua/files/monarch-noctua.conf
[program:monarch-noctua]
command=bash -c "./node_modules/wait-on/bin/wait-on -t 500000 -i 1000 tcp:3400 && ./node_modules/gulp/bin/gulp.js run-noctua"
directory=/opt/noctua
autostart=true
autorestart=true
startretries=3
log_stderr=true
stderr_logfile=/var/log/monarch-noctua.err
stdout_logfile=/var/log/monarch-noctua.out
environment=NODE_PATH="./lib/monarch"
user=TARGET_USER
stopasgroup=true

[eventlistener:httpok]
command=httpok -p monarch-noctua http://localhost:8910/
events=TICK_60
